# ansible/roles/clamav/files/scan_use_clamav.timer
[Unit]
Description=clamav scan all files on boot

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
